# kde/systemd-user-services/update-adlists.service
[Unit]
Description=Update OpenSnitch adlists

[Service]
Type=oneshot
ExecStart=/home/abuturab/.config/scripts/update_adlists.sh
