# kde/systemd-user-services/update-iplists.service
[Unit]
Description=Update OpenSnitch iplists

[Service]
Type=oneshot
ExecStart=/home/abuturab/.config/scripts/update_iplists.sh
