# hyprland/SystemdServices/.config/systemd/user/powerProfileSwitcher.service
[Unit]
Description=Dynamic power profile based on battery state
After=graphical.target

[Service]
ExecStart=%h/.config/scripts/powerProfileSwitcher.sh
Environment=STARTUP_WAIT=30s
Restart=always
RestartSec=5s

# --- Sandboxing ---
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

# Allow read/exec from ~/.local/bin where the script lives
ReadWritePaths=%h/.config/scripts

[Install]
WantedBy=default.target
